# deploy/quadlet/mcp-obsidian-remote.container
[Unit]
Description=OAuth-protected MCP server for Obsidian vaults
After=network-online.target

[Container]
Image=mcp-obsidian-remote.build
ContainerName=mcp-obsidian-remote
UserNS=auto
Volume=/path/to/vault:/data/vault:Z,U
EnvironmentFile=/etc/mcp-obsidian-remote/env
PublishPort=127.0.0.1:3000:3000
RunInit=true
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
